# deploy/systemd/pilingtrack-backup.timer
[Unit]
Description=Nightly PilingTrack Postgres backup

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
