# deploy/systemd/pilingtrack-disk-guard.timer
[Unit]
Description=Run PilingTrack disk guard every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
